# deploy/com.tradebot.paperbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.tradebot.paperbot — runs the whole PAPER bot (intraday ORB+VWAP +
  options short strangle) at 15:45 IST on weekdays (Mon–Fri) via
  run_paper_bot.sh. Simulated only; places no real orders.
  Install:  cp deploy/com.tradebot.paperbot.plist ~/Library/LaunchAgents/
            launchctl load -w ~/Library/LaunchAgents/com.tradebot.paperbot.plist
  Remove:   launchctl unload -w ~/Library/LaunchAgents/com.tradebot.paperbot.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.paperbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/samvid/projects/tradebot/run_paper_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/samvid/projects/tradebot</string>

    <!-- 15:45 local (IST) Mon(1)–Fri(5); launchd runs a missed slot on wake. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/samvid/projects/tradebot/logs/paperbot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/samvid/projects/tradebot/logs/paperbot.log</string>
</dict>
</plist>
