# deploy/com.tradebot.intraday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.tradebot.intraday — auto-runs the intraday ORB PAPER simulator at
  15:45 IST on weekdays (Mon–Fri). Simulated only; places no real orders.
  Install:  cp deploy/com.tradebot.intraday.plist ~/Library/LaunchAgents/
            launchctl load -w ~/Library/LaunchAgents/com.tradebot.intraday.plist
  Remove:   launchctl unload -w ~/Library/LaunchAgents/com.tradebot.intraday.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.intraday</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/samvid/projects/tradebot/.venv/bin/python</string>
        <string>/Users/samvid/projects/tradebot/intraday_sim.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/samvid/projects/tradebot</string>

    <!-- 15:45 local (IST) on Mon(1) Tue(2) Wed(3) Thu(4) Fri(5).
         launchd runs a missed calendar event once when the Mac next wakes. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>45</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/samvid/projects/tradebot/logs/intraday.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/samvid/projects/tradebot/logs/intraday.log</string>
</dict>
</plist>
